# package/snblink/files/sn_blink.init
#!/bin/sh /etc/rc.common
SERVICE_USE_PID=1

START=99

start() {
	sn_blink.sh &
}

stop() {
	killall sn_blink
	killall sn_blink.sh
}

