# package/jd/files/jd.init
#!/bin/sh /etc/rc.common

SERVICE_USE_PID=1
START=97
start() {
	jt &
}
