# package/xl2tpd/files/xl2tpd.init
#!/bin/sh /etc/rc.common
SERVICE_USE_PID=1
START=80

start() {
	ln -s /usr/sbin/pppd /usr/sbin/pppd-l2tp
	mkdir -p /var/run/xl2tpd
}
